# docer/web_agent/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:consul]
command=/bin/agent --config-dir /etc/consul.d
autorestart=true

[program:nginx]
command=/usr/sbin/nginx
autorestart=true
